# deploy/railway/arq-worker.service
[Unit]
Description=Learning Space ARQ Worker (Burst Mode)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/app
ExecStart=/app/.venv/bin/python -m workers.run_worker --burst
Environment=REDIS_URL=rediss://USER:PASSWORD@HOST:PORT
Environment=DATABASE_URL=postgresql+asyncpg://USER:PASSWORD@HOST:PORT/db
Environment=NEO4J_URI=neo4j+s://XXXXX.databases.neo4j.io
Environment=NEO4J_USERNAME=neo4j
Environment=NEO4J_PASSWORD=PASSWORD
Environment=PYTHONPATH=/app
Environment=PATH=/app/.venv/bin:/usr/local/bin:/usr/bin:/bin
StandardOutput=journal
StandardError=journal
User=app
Group=app

[Install]
WantedBy=multi-user.target
